# package/boot/rtk-init/src/android.init
#!/bin/sh /etc/rc.common
# 2020 jjm2473 porting openwrt 19.07+ to zidoo z9s

START=01

boot() {
	# notify android continue
	touch /android/.coldplug_done
	touch /rom/android/.coldplug_done

	# wait wifi
	for i in `seq 1 30`; do
		if [ -d /sys/class/ieee80211/phy0 ]; then
			break
		fi
		sleep 1
	done
}
