# target/linux/realtek/base-files/etc/uci-defaults/00-uhttpd_disable_redirect_https
#!/bin/sh

if [ "`uci -q get uhttpd.main.redirect_https`" == "1" ]; then
	uci set uhttpd.main.redirect_https=0
	uci commit uhttpd
fi

exit 0
