# target/linux/realtek/base-files/etc/uci-defaults/99-unhandle_services
#!/bin/sh

for SERVICE in transmission; do
    uci set "ucitrack.@$SERVICE[-1].exec=/usr/sbin/start_unhandled $SERVICE"
done
